# ais-emissions-matcher.service
[Unit]
Description=AIS Emissions Data Matcher
After=network.target

[Service]
Type=simple
User=erik
WorkingDirectory=/var/www/apihub
ExecStart=/var/www/apihub/venv/bin/python /var/www/apihub/emissions_matcher.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ais-emissions-matcher

[Install]
WantedBy=multi-user.target
