# config/systemd/ais-collector.service
[Unit]
Description=AIS Vessel Data Collector
After=network.target

[Service]
Type=simple
User=erik
WorkingDirectory=/var/www/apihub
Environment="PATH=/var/www/apihub/venv/bin"
ExecStart=/var/www/apihub/venv/bin/python /var/www/apihub/src/collectors/ais_collector.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ais-collector

[Install]
WantedBy=multi-user.target
