# config/systemd/ais-gt-scraper.service
[Unit]
Description=AIS Gross Tonnage Scraper
After=network.target

[Service]
Type=simple
User=erik
WorkingDirectory=/var/www/apihub
ExecStart=/var/www/apihub/venv/bin/python /var/www/apihub/src/services/gt_scraper.py
Restart=always
RestartSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ais-gt-scraper

[Install]
WantedBy=multi-user.target
